# src/MCP/supervisord.conf
[supervisord]
nodaemon=true

[program:weather]
command=python /app/weather_server.py
autostart=true
autorestart=true
stderr_logfile=/var/log/weather.err.log
stdout_logfile=/var/log/weather.out.log

